# systemd/weather.timer
[Unit]
Description=Refresh weather data

[Timer]
OnBootSec=10
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
